# example_config/client/usr/lib/systemd/system/sound.service
[Unit]
Description=Sound

[Service]
ExecStart=/home/pi/bin/optimize_sound

[Install]
WantedBy=multi-user.target
