# example_config/client/lib/systemd/system/mpd.service.d/override.conf
[Unit]
Wants=systemd-networkd-wait-online.service
#Wants=NetworkManager-wait-online.service
After=sound.target local-fs.target network-online.target svscan.service
RequiresMountsFor=/home/pi/MDrive

[Service]
ExecStartPost=/usr/bin/svc -u /service/mpdev /service/fserver /service/upmpdcli
ExecStop=/usr/bin/svc -d /service/mpdev /service/upmpdcli
ExecStopPost=killall upmpdcli

[Install]
Also=systemd-networkd-wait-online.service
#Also=NetworkManager-wait-online.service
